# systemd/rpm-ostreed-automatic.service.d/notify.conf
# Drop-in that pulls rpm-ostree-notify.service into every run of the automatic updater.
[Unit]
Wants=rpm-ostree-notify.service
